# templates/rsnapshot@.service.j2
[Unit]
Description=rsnapshot (%I) backup

[Service]
Type=oneshot
Nice=19
IOSchedulingClass=idle
ExecStart={{ wtd_rsnapshot_service_executeStart }}
